# config/systemd/user/templates/nightlock.service
[Unit]
Description=Lock Screen Automatically at Night

[Service]
Type=oneshot
ExecStart=/usr/bin/xautolock -locknow
